# Gentoo_ExpertBook/cron/cron.weekly/fstrim
#!/bin/sh
# trim all mounted file systems which support it
# /sbin/fstrim --all || true
echo "*** $(date -R) ***" | logger
echo "Starting fstrim..." | logger
/sbin/fstrim -av | logger
