# HomeServer/cron.d/raid-check
# Run system wide raid-check once a week on Sunday at 1am by default
# 0 1 * * Sun root /usr/sbin/raid-check
0 13 * * Mon root /usr/sbin/raid-check

